# deployment/systemd/ap-sso.service
[Unit]
Description=Accounts Payable SSO Frontend
After=network-online.target ap-mock-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/demo_ready_togo/accounts-payable-chatbot
ExecStart=/root/demo_ready_togo/accounts-payable-chatbot/.venv/bin/python /root/demo_ready_togo/accounts-payable-chatbot/mock_sso/main.py
Restart=always
RestartSec=3
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
